# files/etc/services.d/caddy/run
#!/usr/bin/execlineb -P
s6-setuidgid caddy
/usr/bin/caddy --conf /etc/Caddyfile
